# dev/apps.ini
[program:hugo]
command=hugo server --bind=0.0.0.0 --port=80 --verbose --buildDrafts --appendPort=false --disableLiveReload
directory = /app/jrmcc
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=false

[program:consul]
command=consul agent -config-file=/config/consul.json
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=false

[program:grunt]
command=/app/grunt.sh
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=false
